# udhcpc
#!/sbin/openrc-run
dev="$(printf "%s\n" "/sys/class/net/en"*)"
dev="${dev##*/}"
command_background="yes"
command="/bin/busybox udhcpc"
command_args="-i ${dev} -s /etc/udhcpc/default.script"
pidfile="/run/udhcpc.pid"
name="udhcpc"

start_pre() {
        ${command} ${command_args} || {
                command_args="-i "wlan0" -s /etc/udhcpc/default.script"
        }
}

depend() {
        need localmount
        provide net
        after modules
}
